# autostart/user/notebook.service
[Unit]
Description = notebook
[Service]
ExecStart = /usr/bin/jupyter notebook /home/yanpan/NoteBook/
[Install]
WantedBy = default.target
